# playbook/roles/postgresql/templates/etc/profile.d/postgresql.sh.j2
PATH=$PATH:{{ postgresql_bin_directory }}
export PATH
